# foo-worker.conf
start on starting foo
stop on stopping foo
